# crypt-expiry-check/crypt-expiry-check.cron
#!/bin/bash

/usr/sbin/crypt-expiry-check -qa -e me@example.com -f /etc/crypt-expiry.checks
