# logwatch-systemd/logwatch.service
[Unit]
Description=Run logwatch

[Service]
Type=oneshot
ExecStart=/usr/bin/logwatch
